# incubator-STX/incubatorSTX.service
[Unit]
Description=incubatorSTX robot server

[Service]
ExecStart=/pharmbio/labrobots-restserver/start-server-incubatorSTX.sh

[Install]
WantedBy=multi-user.target
